# com.kimaireviews.scheduler.14.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kimaireviews.scheduler.14</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/jungmyungkim/Library/Application Support/Claude/kimaireviews/scheduler.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jungmyungkim/Library/Application Support/Claude/kimaireviews</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>14</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jungmyungkim/Library/Application Support/Claude/kimaireviews/logs/launchd-14.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jungmyungkim/Library/Application Support/Claude/kimaireviews/logs/launchd-14-err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/jungmyungkim</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
